# data/eos-companion-app.service
[Unit]
Description=EOS Companion App Service
After=network.target eos-companion-app.socket
Requires=eos-companion-app.socket

[Service]
ExecStart=/usr/bin/flatpak run --no-desktop com.endlessm.CompanionAppService

# Not ideal, we may want to do this in the user session
User=app-manager

PrivateTmp=yes
